# weatherstation.service
# Example systemd service file for manual installation.
# The automated installer (install.sh) generates its own service file
# with the correct paths for your system.

[Unit]
Description=Weather Station e-Paper Display
After=network.target

[Service]
EnvironmentFile=%h/.env
ExecStart=%h/.local/bin/weatherstation
StandardOutput=append:/var/log/weatherstation.log
StandardError=append:/var/log/weatherstation.log
Restart=always
User=%u
Type=idle
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target
